# OpenGear-ACM500x-4.1.1u2-devkit-20180508/vendors/OpenGear/ACM500x-MG/etc/default/inittab
https:unknown:/bin/stunnel /etc/config/https.conf
sshd:unknown:/bin/sshd -r -D -o AllowUsers=* -o AllowTcpForwarding=yes -o SDT=yes
dial:unknown:/bin/pppd call dialout.console persist nodetach
zebra:unknown:/bin/zebra
ripd:unknown:/bin/ripd
