# OpenGear-ACM500x-4.1.1u2-devkit-20180508/vendors/OpenGear/romfs/etc/inittab
inet:unknown:/bin/inetd
cron:unknown:/usr/sbin/crond -l 9 -S -f
fltd:unknown:/bin/flatfsd
klog:unknown:/sbin/klogd -n
slog:unknown:/sbin/syslogd -n
